# agent/rpm/SOURCES/conf/cron.d/iostat-cron.conf
*/1 * * * * root /usr/lib/zabbix/externalscripts/iostat-cron.sh
