# agent/rpm/SOURCES/conf/cron.d/multicpu-cron.conf
*/1 * * * * root /usr/lib/zabbix/externalscripts/multicpu-cron.sh
